# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
user = root
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisor/supervisord.log
childlogdir = /var/log/supervisor

[program:dnsmasq]
command = dnsmasq -k -C /etc/dnsmasq.conf
autostart = true
autorestart = false
stdout_logfile=/var/log/supervisor/dnsmasq.out
stderr_logfile=/var/log/supervisor/dnsmasq.err
